# docker/supervisor/forum_version_supr.conf
[program:set_forum_image_version]
command = /opt/docker/bin/service.d/set_forum_version.sh
startsecs = 10
autorestart = false
startretries = 1
