# docker/supervisor/forum_version_supr.conf
[program:set_forum_image_version]
command = /opt/docker/bin/service.d/set_forum_version.sh
startsecs = 0
autorestart = unexpected
exitcodes = 0
startretries = 3
priority = 2000
